# crontab
#https://crontab.guru/examples.html
0 */6 * * * /root/imhere2.sh
